# services/celery/celery.conf
# if changed, this must be manually copied to /etc/init/celery.conf
description "celery"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
exec sudo -u vagrant /home/vagrant/.virtualenvs/perma/bin/celery -A perma worker --loglevel=info -Q celery --logfile=/vagrant/services/logs/celery.log --workdir=/vagrant/perma_web/

